# solis-management/solis-management.service
[Unit]
Description=Solis Inverter Management API
After=network.target
Documentation=file:/media/grafanas-vm/git/solar-management/solis-management/README.md

[Service]
Type=simple
User=ilia
WorkingDirectory=/media/grafanas-vm/git/solar-management/solis-management
ExecStart=/usr/bin/python3 /media/grafanas-vm/git/solar-management/solis-management/solis-api.py --host 0.0.0.0 --port 5000
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=solis-management

[Install]
WantedBy=multi-user.target
